# kubernetes/supervisord.conf
[program:node]
command=node /usr/src/app/dist/server
stopsignal=KILL
autostart=true
autorestart=true
startsecs=20
